# ground-station/config/direwolf.service
# NOTE: Direwolf is now managed by umbgs directly (rtl_fm | direwolf pipeline).
# This service file is kept only as a stub for compatibility.
# Do NOT enable this service — umbgs handles direwolf lifecycle.

[Unit]
Description=Direwolf APRS TNC (legacy stub — managed by umbgs)

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
