# ground-station/config/umbgs-watchdog.service
[Unit]
Description=UMB Ground Station Watchdog
After=umbgs.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/umbgs-watchdog.sh
